# systemd/fuzzyclock.service
[Unit]
Description=E-Ink Fuzzy Clock Daemon
After=network.target

[Service]
# __REPO_DIR__ and __USER__ are substituted by deploy.sh at install time.
ExecStart=/usr/bin/python3 __REPO_DIR__/fuzzyclock_daemon.py
Restart=always
User=__USER__
WorkingDirectory=__REPO_DIR__
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
